# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log

# ── PostgreSQL (starts first) ──────────────────────────────────────────────────
[program:postgres]
command=PG_BIN_PLACEHOLDER/postgres -D /var/lib/postgresql/data
user=postgres
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ── Uvicorn (wait-for-postgres.sh polls pg_isready before launching) ───────────
[program:uvicorn]
command=/wait-for-postgres.sh
directory=/app
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ── Nginx ─────────────────────────────────────────────────────────────────────
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ── Telegram bot (optional — exits 0 cleanly when BOT_TOKEN is not set) ───────
[program:telegram-bot]
command=/start-telegram-bot.sh
autostart=true
autorestart=unexpected
exitcodes=0
startsecs=0
priority=40
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
